# config/btsync.upstart.conf
description "NGCDN BTSync Service"
author "matt@mrneel.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

script
	exec btsync --config /etc/btsync.conf
end script